# events/EC2/StartInstances.json
{
    "eventName": "StartInstances",
    "eventSource": "ec2.amazonaws.com",
    "awsService": "EC2",
    "description": "Starts an Amazon EBS-backed instance that you've previously stopped.",
    "mitreAttackTactics": [
        "TA0003 - Persistence",
        "TA0040 - Impact"
    ],
    "mitreAttackTechniques": [
        "T1098 - Account Manipulation",
        "T1496 - Resource Hijacking"
    ],
    "mitreAttackSubTechniques": [],
    "unverifiedMitreAttackTechniques": [
        {
            "technique": "T1036 - Masquerading",
            "reason": "Adversaries could rename stopped instances to appear legitimate and start them without raising alarms."
        },
        {
            "technique": "T1053 - Scheduled Task/Job",
            "reason": "Attackers might schedule tasks to automatically start stopped instances at certain times to execute malicious actions"
        },
        {
            "technique": "T1105 - Ingress Tool Transfer",
            "reason": "Once the instance is started, adversaries could transfer tools and malware to the instance for execution"
        },
        {
            "technique": "T1219 - Remote Access Software",
            "reason": "Adversaries might start instances that have remote access tools installed to regain control over the environment."
        },
        {
            "technique": "T1562 - Impair Defenses",
            "reason": "Starting instances can impair defenses by creating new workloads that may not be monitored by existing security tools, enabling attackers to perform malicious activities without detection."
        },
        {
            "technique": "T1578 - Modify Cloud Compute Infrastructure",
            "reason": "Attackers can directly use the StartInstances API call to manipulate the state of instances, aiding in persistence and execution of tasks."
        },
        {
            "technique": "T1074 - Data Staged",
            "reason": "Starting an instance can be used to stage data locally before exfiltration."
        }
    ],
    "usedInWild": true,
    "incidents": [
        {
            "description": "When a Zero Day and Access Keys Collide in the Cloud: Responding to the SugarCRM Zero-Day Vulnerability",
            "link": "https://unit42.paloaltonetworks.com/sugarcrm-cloud-incident-black-hat/"
        }
    ],
    "researchLinks": [
        {
            "description": "Executing commands through EC2 user data",
            "link": "https://securitylabs.datadoghq.com/cloud-security-atlas/attacks/executing-commands-through-user-data/"
        },
        {
            "description": "Attack Paths Into VMs in the Cloud",
            "link": "https://unit42.paloaltonetworks.com/cloud-virtual-machine-attack-vectors/"
        }
    ],
    "securityImplications": "Attackers might use StartInstances to reactivate dormant EC2 instances or after having modified the user data for execution of commands.",
    "alerting": [],
    "simulation": [
        {
            "type": "commandLine",
            "value": "aws ec2 start-instances --instance-ids TrailDiscoverInstanceID"
        }
    ],
    "permissions": "https://aws.permissions.cloud/iam/ec2#ec2-StartInstances"
}
