# events/ECS/RegisterTaskDefinition.json
{
    "eventName": "RegisterTaskDefinition",
    "eventSource": "ecs.amazonaws.com",
    "awsService": "ECS",
    "description": "Registers a new task definition from the supplied family and containerDefinitions.",
    "mitreAttackTactics": [
        "TA0040 - Impact"
    ],
    "mitreAttackTechniques": [
        "T1496 - Resource Hijacking"
    ],
    "mitreAttackSubTechniques": [],
    "unverifiedMitreAttackTechniques": [
        {
            "technique": "T1053 - Scheduled Task/Job",
            "reason": "Registering a task definition can be leveraged to create scheduled tasks within ECS, allowing for persistence and automated execution of malicious tasks."
        },
        {
            "technique": "T1105 - Ingress Tool Transfer",
            "reason": "Task definitions could be used to download and execute additional tools or scripts from external sources"
        },
        {
            "technique": "T1562 - Impair Defenses",
            "reason": "ECS tasks can be configured to disable or modify security tools within the container environment, aiding in defense evasion."
        },
        {
            "technique": "T1059 - Command and Scripting Interpreter",
            "reason": "The task definitions can contain Unix shell commands, facilitating execution of malicious scripts or commands."
        },
        {
            "technique": "T1036 - Masquerading",
            "reason": "Malicious task definitions can be disguised as legitimate ones to evade detection and blend in with normal operations"
        },
        {
            "technique": "T1027 - Obfuscated Files or Information",
            "reason": "Container definitions within ECS can include obfuscated or packed scripts and binaries, making detection harder."
        }
    ],
    "usedInWild": true,
    "incidents": [
        {
            "description": "Tales from the cloud trenches: Amazon ECS is the new EC2 for crypto mining",
            "link": "https://securitylabs.datadoghq.com/articles/tales-from-the-cloud-trenches-ecs-crypto-mining/"
        },
        {
            "description": "GuardDuty Extended Threat Detection uncovers cryptomining campaign on Amazon EC2 and Amazon ECS",
            "link": "https://aws.amazon.com/blogs/security/cryptomining-campaign-targeting-amazon-ec2-and-amazon-ecs/"
        }
    ],
    "researchLinks": [],
    "securityImplications": "Attackers might use RegisterTaskDefinition to deploy containers with malicious tasks in AWS ECS.",
    "alerting": [
        {
            "type": "sigma",
            "value": "https://github.com/SigmaHQ/sigma/blob/master/rules/cloud/aws/cloudtrail/aws_ecs_task_definition_cred_endpoint_query.yml"
        }
    ],
    "simulation": [
        {
            "type": "commandLine",
            "value": "aws ecs register-task-definition --family 'xtdb-bench-dev' --network-mode 'awsvpc' --container-definitions '[{\"name\":\"bench-container\", \"cpu\":2048, \"memory\":4092 }]'"
        }
    ],
    "permissions": "https://aws.permissions.cloud/iam/ecs#ecs-RegisterTaskDefinition"
}
