# events/SSM/SendCommand.json
{
    "eventName": "SendCommand",
    "eventSource": "ssm.amazonaws.com",
    "awsService": "SSM",
    "description": "Runs commands on one or more managed nodes.",
    "mitreAttackTactics": [
        "TA0008 - Lateral Movement",
        "TA0002 - Execution"
    ],
    "mitreAttackTechniques": [
        "T1021 - Remote Services",
        "T1651 - Cloud Administration Command"
    ],
    "mitreAttackSubTechniques": [

    ],
    "unverifiedMitreAttackTechniques": [
        {
            "technique": "T1059 - Command and Scripting Interpreter",
            "reason": "The SendCommand API is designed to execute commands on managed nodes, which directly involves the use of command and scripting interpreters to run scripts or commands."
        },
        {
            "technique": "T1053 - Scheduled Task/Job",
            "reason": "The SendCommand API can be used to create or modify scheduled tasks on managed nodes, enabling the execution of commands at specified times, which is essential for maintaining persistence."
        },
        {
            "technique": "T1105 - Ingress Tool Transfer",
            "reason": "Attackers can use SendCommand to download and execute additional tools or payloads on the managed nodes, which is directly relevant to executing commands that facilitate further compromise."
        },
        {
            "technique": "T1569 - System Services",
            "reason": "The SendCommand API can start, stop, or restart system services on managed nodes, allowing for the execution of commands that may serve various purposes, including persistence or privilege escalation."
        }
    ],
    "usedInWild": true,
    "incidents": [
        {
            "description": "Navigating the Cloud: Exploring Lateral Movement Techniques",
            "link": "https://unit42.paloaltonetworks.com/cloud-lateral-movement-techniques/"
        }
    ],
    "researchLinks": [
        {
            "description": "Hunting for signs of persistence in the cloud: an IR guide following the CircleCI incident",
            "link": "https://www.wiz.io/blog/hunting-for-signs-of-persistence-in-the-cloud-an-ir-guide"
        },
        {
            "description": "Run Shell Commands on EC2 with Send Command or Session Manager",
            "link": "https://hackingthe.cloud/aws/post_exploitation/run_shell_commands_on_ec2/"
        },
        {
            "description": "Attack Paths Into VMs in the Cloud",
            "link": "https://unit42.paloaltonetworks.com/cloud-virtual-machine-attack-vectors/"
        }
    ],
    "securityImplications": "Attackers might use SendCommand to execute malicious commands on managed instances.",
    "alerting": [],
    "simulation": [
        {
            "type": "commandLine",
            "value": "aws ssm send-command --instance-ids \"TrailDiscoverInstanceID\" --document-name \"AWS-RunShellScript\" --parameters commands=ls --output text"
        },
        {
            "type": "stratusRedTeam",
            "value": "https://stratus-red-team.cloud/attack-techniques/AWS/aws.execution.ssm-send-command"
        }
    ],
    "permissions": "https://aws.permissions.cloud/iam/ssm#ssm-SendCommand"
}