# events/EventBridge/PutRule.json
{
    "eventName": "PutRule",
    "eventSource": "events.amazonaws.com",
    "awsService": "events",
    "description": "Creates or updates the specified rule.",
    "mitreAttackTactics": [
        "TA0040 - Impact",
        "TA0005 - Defense Evasion",
        "TA0003 - Persistence"
    ],
    "mitreAttackTechniques": [
        "T1489 - Service Stop",
        "T1578 - Modify Cloud Compute Infrastructure",
        "T1546 - Event Triggered Execution"
    ],
    "mitreAttackSubTechniques": [],
    "unverifiedMitreAttackTechniques": [
        {
            "technique": "T1205 - Traffic Signaling",
            "reason": "EventBridge rules can be configured to trigger signals that facilitate command and control communication, masking malicious traffic as legitimate event triggers."
        },
        {
            "technique": "T1053 - Scheduled Task/Job: Scheduled Task",
            "reason": "Creating or updating EventBridge rules can schedule tasks or jobs that perform malicious activities without user intervention."
        },
        {
            "technique": "T1070 - Indicator Removal",
            "reason": "By manipulating EventBridge rules, attackers can potentially alter the flow of logs and events to hide their activities."
        },
        {
            "technique": "T1562 - Impair Defenses",
            "reason": "By updating EventBridge rules, attackers can disable or modify security tools and alerts, impairing defenses and ensuring continued access."
        }
    ],
    "usedInWild": true,
    "incidents": [
        {
            "description": "Cloud Security Stories: From Risky Permissions to Ransomware Execution",
            "link": "https://sonraisecurity.com/blog/cloud-security-stories-from-risky-permissions-to-ransomware-execution/"
        }
    ],
    "researchLinks": [
        {
            "description": "Modify GuardDuty Configuration",
            "link": "https://hackingthe.cloud/aws/avoiding-detection/modify-guardduty-config/"
        }
    ],
    "securityImplications": "Attackers might use PutRule in AWS EventBridge to create unauthorized event rules, potentially automating malicious actions to gain persistence or triggering unwarranted responses within the environment.",
    "alerting": [],
    "simulation": [
        {
            "type": "commandLine",
            "value": "aws events put-rule --name TrailDiscoverRule --schedule-expression 'rate(5 minutes)' --state ENABLED --description \"TrailDiscover rule\""
        }
    ],
    "permissions": "https://aws.permissions.cloud/iam/events#events-PutRule"
}
