# events/Lambda/CreateFunction20150331.json
{
    "eventName": "CreateFunction20150331",
    "eventSource": "lambda.amazonaws.com",
    "awsService": "Lambda",
    "description": "Creates a Lambda function.",
    "mitreAttackTactics": [
        "TA0003 - Persistence",
        "TA0004 - Privilege Escalation",
        "TA0040 - Impact"
    ],
    "mitreAttackTechniques": [
        "T1098 - Account Manipulation",
        "T1496 - Resource Hijacking"
    ],
    "mitreAttackSubTechniques": [],
    "unverifiedMitreAttackTechniques": [
        {
            "technique": "T1059 - Command and Scripting Interpreter",
            "reason": "The Lambda function can be configured to execute JavaScript code, enabling attackers to run malicious scripts."
        },
        {
            "technique": "T1070 - Indicator Removal",
            "reason": " By using Lambda, attackers can delete logs or files to evade detection."
        },
        {
            "technique": "T1562 - Impair Defenses",
            "reason": "Attackers might create Lambda functions designed to disable security monitoring tools or alerts."
        },
        {
            "technique": "T1071 - Application Layer Protocol",
            "reason": "Lambda functions can communicate over standard web protocols, enabling Command and Control communication that blends with regular traffic."
        },
        {
            "technique": "T1027 - Obfuscated Files or Information",
            "reason": "If attackers use code-signing configurations and different deployment packaging (e.g., obfuscated container images or encrypted .zip archives), it can help evade detection by concealing the true nature of the function code."
        },
        {
            "technique": "T1053 - Scheduled Task/Job",
            "reason": "Attackers might schedule Lambda functions to execute at specific intervals, providing a means of persistence or delayed execution."
        }
    ],
    "usedInWild": true,
    "incidents": [
        {
            "description": "Mining Crypto",
            "link": "https://twitter.com/jonnyplatt/status/1471453527390277638"
        },
        {
            "description": "Cloud Security Stories: From Risky Permissions to Ransomware Execution",
            "link": "https://sonraisecurity.com/blog/cloud-security-stories-from-risky-permissions-to-ransomware-execution/"
        },
        {
            "description": "Leaked Environment Variables Allow Large-Scale Extortion Operation in Cloud Environments",
            "link": "https://unit42.paloaltonetworks.com/large-scale-cloud-extortion-operation/"
        },
        {
            "description": "Tales from the cloud trenches: The Attacker doth persist too much, methinks",
            "link": "https://securitylabs.datadoghq.com/articles/tales-from-the-cloud-trenches-the-attacker-doth-persist-too-much/"
        },
        {
            "description": "GuardDuty Extended Threat Detection uncovers cryptomining campaign on Amazon EC2 and Amazon ECS",
            "link": "https://aws.amazon.com/blogs/security/cryptomining-campaign-targeting-amazon-ec2-and-amazon-ecs/"
        }
    ],
    "researchLinks": [
        {
            "description": "Hunting for signs of persistence in the cloud: an IR guide following the CircleCI incident",
            "link": "https://www.wiz.io/blog/hunting-for-signs-of-persistence-in-the-cloud-an-ir-guide"
        },
        {
            "description": "AWS IAM Privilege Escalation Techniques",
            "link": "https://hackingthe.cloud/aws/exploitation/iam_privilege_escalation/"
        }
    ],
    "securityImplications": "Attackers might use CreateFunction to deploy malicious code or functions, depending on the scenario this might allow the attacker to gain persistence, escalate privileges, or hijack resources.",
    "alerting": [],
    "simulation": [
        {
            "type": "commandLine",
            "value": "aws lambda create-function --function-name my-function --runtime nodejs18.x --code S3Bucket=string --role arn:aws:iam::123456789012:role/service-role/MyTestFunction-role-tges6bf4"
        }
    ],
    "permissions": "https://aws.permissions.cloud/iam/lambda#lambda-CreateFunction"
}
